# config/crontab/crontab
* * * * * /var/www/html/jacktrade/bin/crunz schedule:run

